# roles/nextcloud/files/nextcloud.service
# /etc/systemd/system
[Unit]
Description=Enable docker/nextcloud on startup

[Service]
Type=simple
ExecStart=/usr/bin/docker-compose --env-file /srv/git/dockers/nextcloud/.env --file /srv/git/dockers/nextcloud/docker-compose.yml up --detach
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
